# etc/supervisord.d/swagger-editor.ini
[program:swagger-editor]
directory=/opt/
command=http-server swagger-editor
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
